# Apps/NexusMac/LaunchAgents/com.kacperpietrzyk.nexus.meetings-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kacperpietrzyk.nexus.meetings-helper</string>
    <key>MachServices</key>
    <dict>
        <key>group.com.kacperpietrzyk.Nexus.meetings-helper</key>
        <true/>
    </dict>
    <key>BundleProgram</key>
    <string>Contents/Library/LoginItems/NexusMeetingsHelper.app/Contents/MacOS/NexusMeetingsHelper</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
